# docker/etc/s6-overlay/services/frontend/finish
#!/usr/bin/with-contenv bash
s6-svc -d /var/run/s6/services/nginx
