# docker/etc/s6-overlay/services/frontend/run
#!/usr/bin/with-contenv bash
cd /app/frontend
exec node_modules/.bin/next start
