# docker/etc/s6-overlay/services/nginx/run
#!/usr/bin/with-contenv bash
exec nginx -g "daemon off;"
